# supervisor_celery.conf
; DataViz Studio - Celery Worker Configuration
; Environment: development
; Generated by deploy.sh

[program:celery_worker]
command=celery -A utils.tasks worker --loglevel=info --concurrency=4 --queues=celery,high_priority,low_priority
directory=/app/backend
user=root
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/dataviz/celery_worker_%(process_num)02d.log
stderr_logfile=/var/log/dataviz/celery_worker_%(process_num)02d_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopasgroup=true
killasgroup=true
priority=998
environment=PYTHONPATH="/app/backend"

[program:celery_beat]
command=celery -A utils.tasks beat --loglevel=info
directory=/app/backend
user=root
numprocs=1
stdout_logfile=/var/log/dataviz/celery_beat.log
stderr_logfile=/var/log/dataviz/celery_beat_err.log
autostart=true
autorestart=true
startsecs=10
priority=999
environment=PYTHONPATH="/app/backend"

[group:celery]
programs=celery_worker,celery_beat
